# media_server/ansible/roles/torrent/files/lib/systemd/system/deluge-web.service
[Unit]
Description=Deluge Bittorrent Web Daemon
After=network-online.target

[Service]
User=milhouse
Group=milhouse
Type=simple
ExecStart=/usr/bin/deluge-web
Restart=on-failure

[Install]
WantedBy=multi-user.target
